# Scripts/RunUO.conf
description "RunUO"
author "Jacob Harrison"
start on runlevel [2345]

exec mono ~/project-era-uo/RunUO/RunUO.exe
